# supervisord.conf

[program:dwc_services]
command=java %{ENV_JAVA_OPTS}s -jar jetty.jar --path %{ENV_CONTEXT}s dwc-services.war
stdout_logfile=syslog
stderr_logfile=syslog

