# supervisord.conf

[program:dwc_services]
command=/root/start.sh
stdout_logfile=syslog
stderr_logfile=syslog

