# supervisord.conf

[program:dwc_services]
command=java %(ENV_JAVA_OPTS)s -jar /root/jetty.jar --path %(ENV_CONTEXT)s /root/dwc-services.war
stdout_logfile=syslog
stderr_logfile=syslog

